# roles/baseline/files/dnf-update.timer
[Unit]
Description=dnf update -y

[Timer]
OnCalendar=Mon..Fri 02:00
AccuracySec=1m

