# roles/baseline/files/dnf-update.service
[Unit]
Description=dnf update -y

[Service]
Type=oneshot
ExecStart=/usr/bin/dnf update -y

